# assets/06_complaint_processor/conf/supervisor/complaint-worker.conf
[program:complaint-worker]
command=/opt/complaint-worker/venv/bin/celery -A tasks worker --loglevel=info --concurrency=3 -Q file_convert
directory=/opt/complaint-worker
user=complaint-worker
autostart=true
autorestart=true
stderr_logfile=/opt/complaint-worker/logs/worker_err.log
stdout_logfile=/opt/complaint-worker/logs/worker.log
environment=
    PATH="/opt/complaint-worker/venv/bin:/usr/local/bin:/usr/bin:/bin",
    HOME="/home/complaint-worker"

; [취약 설정] No sandboxing options configured
; 올바른 설정: chroot, seccomp, network namespace isolation should be applied
