# assets/C14_데이터수집서버/conf/systemd/datacollector.service
# /etc/systemd/system/datacollector.service
# C14 데이터 수집·관리 서버 — FastAPI 서비스

[Unit]
Description=C4I Data Collector API Server
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/datacollector/app
ExecStart=/opt/datacollector/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000 --log-level info --access-log
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
