# assets/06_민원처리서버/conf/systemd/complaint-worker.service
# =============================================================================
# systemd 유닛 — Celery 워커 (Supervisor 경유)
# 민원 처리 서버 (자산 06)
# =============================================================================
# [취약점] VULN-06-05: 보안 샌드박싱 옵션 미적용
# 올바른 구현: ProtectSystem=strict, PrivateTmp=true, NoNewPrivileges=true 등

[Unit]
Description=Complaint Processing Worker (Supervisor)
After=network.target redis.service
Wants=redis.service

[Service]
Type=forking

# Supervisor를 통해 Celery 워커 관리
ExecStart=/usr/bin/supervisord -c /etc/supervisor/supervisord.conf
ExecStop=/usr/bin/supervisorctl shutdown
ExecReload=/usr/bin/supervisorctl reload

# [취약점] VULN-06-05: complaint-worker 사용자로 실행하지만
# chroot, seccomp, capabilities 제한 없음
# 올바른 구현:
# ProtectSystem=strict
# ProtectHome=true
# PrivateTmp=true
# NoNewPrivileges=true
# CapabilityBoundingSet=
# SystemCallFilter=@system-service

User=root
Group=root

Restart=on-failure
RestartSec=10

# 로그
StandardOutput=journal
StandardError=journal
SyslogIdentifier=complaint-worker

# 환경 파일
EnvironmentFile=/opt/complaint-worker/.env

[Install]
WantedBy=multi-user.target
