# assets/06_민원처리서버/conf/systemd/redis.service
# =============================================================================
# systemd 오버라이드 — Redis 서버
# 민원 처리 서버 (자산 06)
# =============================================================================
# /etc/systemd/system/redis-server.service.d/override.conf 으로 배포

[Service]
# 커스텀 설정 파일 사용
ExecStart=
ExecStart=/usr/bin/redis-server /etc/redis/redis.conf

# [취약 설정] VULN-06-03: Redis가 0.0.0.0에 바인딩되므로
# ProtectNetwork 등 보안 옵션을 적용하지 않음
# 올바른 구현: PrivateNetwork=true 또는 IPAddressAllow= 로 접근 제한
