# assets/03_민원접수서버/conf/systemd/minio.service
[Unit]
Description=MinIO Object Storage
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=minio-user
Group=minio-user
EnvironmentFile=/opt/minio/minio.env
ExecStart=/opt/minio/minio server /var/lib/minio --console-address ":9001"
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
