# assets/C15_상황요약AI/conf/systemd/summary-timer.timer
# /etc/systemd/system/summary-gen.timer
# C15 상황 요약 AI — 15분 주기 브리핑 생성 타이머

[Unit]
Description=C4I Summary Generation Timer (15min)

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
AccuracySec=30s

[Install]
WantedBy=timers.target
