# assets/C15_상황요약AI/conf/systemd/summary-timer.service
# /etc/systemd/system/summary-gen.service
# C15 상황 요약 AI — 브리핑 생성 서비스

[Unit]
Description=C4I Summary Generation Service
After=network-online.target ollama.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/summary-ai/scripts
ExecStart=/opt/summary-ai/venv/bin/python /opt/summary-ai/scripts/summary_pipeline.py
Environment=PYTHONUNBUFFERED=1
Environment=C4I_API_URL=http://192.168.130.12:8000
Environment=C4I_API_KEY=dev-key-12345
Environment=SUMMARY_OLLAMA_URL=http://localhost:11434
