# assets/06_민원처리서버/conf/supervisor/celery-worker.conf
; =============================================================================
; Supervisor 설정 — Celery 워커 프로세스
; 민원 처리 서버 (자산 06)
; =============================================================================
; [취약점] VULN-06-05: 샌드박싱 없이 서비스 사용자 권한으로 실행
; 올바른 구현: chroot, seccomp, capabilities 제한, 네트워크 네임스페이스 격리

[program:complaint-worker]
command=/opt/complaint-worker/venv/bin/celery -A celery_app worker
    --loglevel=info
    --concurrency=3
    -Q file_convert
    --without-heartbeat

directory=/opt/complaint-worker
user=complaint-worker

; 자동 시작/재시작
autostart=true
autorestart=true
startsecs=10
startretries=3

; 종료 시그널
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true

; 로그
stderr_logfile=/opt/complaint-worker/logs/worker_err.log
stdout_logfile=/opt/complaint-worker/logs/worker.log
stderr_logfile_maxbytes=50MB
stdout_logfile_maxbytes=50MB
stderr_logfile_backups=5
stdout_logfile_backups=5

; 환경변수
; [취약점] VULN-06-05: 워커 환경 그대로 상속, .env 경로 직접 노출
; 올바른 구현: 최소한의 환경변수만 전달, 시크릿은 vault에서 로딩
environment=
    PATH="/opt/complaint-worker/venv/bin:/usr/local/bin:/usr/bin:/bin",
    HOME="/home/complaint-worker",
    LANG="en_US.UTF-8"

; 프로세스 수
numprocs=1
process_name=%(program_name)s

; 우선순위
priority=999
