# assets/13_SCADA서버/conf/systemd/nodered.service
# Asset #13 SCADA 서버 - Node-RED systemd unit
# 경로: /etc/systemd/system/nodered.service
# PLC 데이터 브릿지 및 대시보드

[Unit]
Description=Node-RED - OT Data Bridge
Documentation=https://nodered.org/docs/
After=network.target tomcat.service

[Service]
Type=simple
User=nodered
Group=nodered

WorkingDirectory=/home/nodered/.node-red
Environment="NODE_RED_HOME=/home/nodered/.node-red"

ExecStart=/usr/bin/node-red --settings /home/nodered/.node-red/settings.js --userDir /home/nodered/.node-red
ExecStop=/bin/kill -SIGTERM $MAINPID

Restart=on-failure
RestartSec=10

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
