# assets/04_내부업무포털/conf/systemd/intranet-portal.service
[Unit]
Description=MOIS Intranet Portal (Gunicorn + Django)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/intranet-portal/src/intranet_portal
EnvironmentFile=/opt/intranet-portal/.env
ExecStart=/opt/intranet-portal/venv/bin/gunicorn \
    intranet_portal.wsgi:application \
    --workers 4 \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --access-logfile /var/log/intranet-portal/gunicorn-access.log \
    --error-logfile /var/log/intranet-portal/gunicorn-error.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
