# assets/14_PLC시뮬레이터/conf/systemd/plc-simulator.service
[Unit]
Description=PLC Simulator - Virtual Sensor Data Generator
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/plc-simulator
Environment="PATH=/opt/plc-simulator/venv/bin:/usr/bin"
ExecStart=/opt/plc-simulator/venv/bin/gunicorn \
    --bind 0.0.0.0:5000 \
    --workers 2 \
    --threads 4 \
    --timeout 120 \
    --access-logfile /var/log/plc-simulator/access.log \
    --error-logfile /var/log/plc-simulator/error.log \
    "plc_simulator:app"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
